# auto-deploy.timer
[Unit]
Description=Auto Deploy Timer — poll Gitee every 60 seconds
Requires=auto-deploy.service

[Timer]
OnBootSec=30
OnUnitActiveSec=60
AccuracySec=5
Unit=auto-deploy.service

[Install]
WantedBy=timers.target
